# extras/glusterrestd.service.in
[Unit]
Description=Gluster REST APIs
After=syslog.target network.target

[Service]
Type=simple
ExecReload=/bin/kill -SIGHUP $MAINPID
KillMode=control-group
ExecStart=/usr/bin/gunicorn -c @SYSCONF_DIR@/glusterrest/gunicorn_config.py main:app --pid /var/run/glusterrest.pid --chdir @GLUSTER_LIBEXECDIR@/glusterrest
PIDFile=/var/run/glusterrestd.pid

[Install]
WantedBy=multi-user.target
